# pkg/scripts/preremove.sh
#!/usr/bin/env bash

systemctl stop lumeweb-relay.service
systemctl disable lumeweb-relay.service
