# pkg/scripts/preremove.sh
#!/usr/bin/env bash

systemctl disable lumeweb-relay.service
